# ansible/roles/girder/templates/girder.service.j2
[Unit]
Description=Girder server
After=network.target

[Service]
User={{  girder_exec_user }}
ExecStart=/usr/local/bin/girder-server
StandardOutput=null
StandardError=null
